# deploy/ifohack.service
# symlink this to /etc/systemd/system/ on a server with systemd
[Unit]
Description=ifohack.live Apps Services
After=network.target

[Service]
Type=simple
Restart=always
ExecStart=$HOME/ifohack2023/deploy/run.sh

[Install]
WantedBy=multi-user.target
